# src/inittab
# This is an example /etc/inittab file for "busybox init" which starts
# busyrc on boot, creates TTY's and handles the killing of services on
# shutdown.

# It is safest to use absolute paths, as PATH doesn't have a set-in-stone value
# It is likely to be: PATH=/sbin:/usr/sbin:/bin:/usr/bin
# But could change by busybox implementation of if another init exports PATH
# before starting busybox

# Start "rc init" on boot
::sysinit:/sbin/rc init

# Setup login shells
# You could use "-/bin/sh" for a direct login shell
# or "agetty -a root tty1 linux" for autologin
# Login shells are normally launched on TTYs 1-6
# A graphical display manager is traditionally launched on tty7
tty1::respawn:agetty tty1 linux
tty2::respawn:agetty tty2 linux
tty3::respawn:agetty tty3 linux
tty4::respawn:agetty tty4 linux
tty5::respawn:agetty tty5 linux
tty6::respawn:agetty tty6 linux
tty7::respawn:agetty tty7 linux
#tty7::respawn:agetty -a your_username tty7 linux

# Shutdown when pressing CTRL+ALT+DEL (disabled by default)
#::ctrlaltdel:kill -USR2 1

# Stop all services on shutdown
::shutdown:rc shutdown

# Killing everything on shutdown
::shutdown:echo :: sending SIGTERM to all
::shutdown:kill -s TERM -1
::shutdown:sleep 1
::shutdown:echo :: sending SIGKILL to all
::shutdown:kill -s KILL -1

# Unmount everything on shutdown
::shutdown:echo :: unmounting everything
::shutdown:umount -a -r
::shutdown:mount -o remount,ro /

